# autonomous-power-profiles/etc_udev_rules.d/61-onacpower.rules
# Rule for when switching to AC Power
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/share/power-profiles/performance.sh"
